# linkat-18.04/linkat-caf-fix-18.04/debian/postinst
#!/bin/bash
# postinst script for linkat-caf-fix
#
# see: dh_installdeb(1)

set -e

# summary of how this script can be called:
#        * <postinst> `configure' <most-recently-configured-version>
#        * <old-postinst> `abort-upgrade' <new version>
#        * <conflictor's-postinst> `abort-remove' `in-favour' <package>
#          <new-version>
#        * <postinst> `abort-remove'
#        * <deconfigured's-postinst> `abort-deconfigure' `in-favour'
#          <failed-install-package> <version> `removing'
#          <conflicting-package> <version>
# for details, see http://www.debian.org/doc/debian-policy/ or
# the debian-policy package

case "$1" in
	configure)
	if [ -f /etc/init.d/CA-DSM ]; then
		if [ -z "$(grep caf\ kill\ all /etc/init.d/CA-DSM)" ]; then
			cp /etc/init.d/CA-DSM /etc/init.d/CA-DSM.bak
			sed -i 's/sd_jexec\ unit/\#sd_jexec\ unit/g' /etc/init.d/CA-DSM
			sed -i '/#sd_jexec\ unit=\.\ shutdown\ >>/a\\t\t\tcaf\ kill\ all\ \>\/dev\/null 2>&1' /etc/init.d/CA-DSM 
		fi
	fi
	if [  -e /usr/share/linkat/linkat-caf-fix/linkat-caf-fix.sh ]; then 
		chmod 750 /usr/share/linkat/linkat-caf-fix/linkat-caf-fix.sh
	fi
	if [ -e /usr/share/linkat/linkat-caf-fix/linkat-caf-cron ]; then
		cp /usr/share/linkat/linkat-caf-fix/linkat-caf-cron /etc/cron.d/.
	fi
	invoke-rc.d cron reload >/dev/null 2>&1
	if [ -e /usr/share/linkat/linkat-caf-fix/scala_centro_CA.csv ]; then
		chmod 640 /usr/share/linkat/linkat-caf-fix/scala_centro_CA.csv
	fi
	if [ -e /usr/share/linkat/linkat-caf-fix/config-ca.sh ]; then
		chmod 755 /usr/share/linkat/linkat-caf-fix/config-ca.sh
	fi	
	if [ -d /opt/configuracio ]; then
		cp /usr/share/linkat/linkat-caf-fix/config-ca.sh /opt/configuracio/config-ca.sh
		chmod 755 /opt/configuracio/config-ca.sh
		chown root:root /opt/configuracio/config-ca.sh
	fi
	if [ -n "$(grep mnt-recuperacio /etc/fstab)" ]; then
		if [ -z "$(mount |grep mnt-recuperacio)" ]; then
			mount /opt/configuracio/mnt-recuperacio
		fi
		if [ -n "$(mount |grep mnt-recuperacio)" ]; then
			if [ -e /opt/configuracio/mnt-recuperacio/config ]; then
				cp /usr/share/linkat/linkat-caf-fix/scala_centro_CA.csv /opt/configuracio/mnt-recuperacio/config/.
				dos2unix /opt/configuracio/mnt-recuperacio/config/scala_centro_CA.csv >/dev/null 2>&1
				chmod 640 /opt/configuracio/mnt-recuperacio/config/scala_centro_CA.csv
			fi
		fi
	fi
	sleep 1
	umount /opt/configuracio/mnt-recuperacio
    ;;

    abort-upgrade|abort-remove|abort-deconfigure)
    ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0
