# linkat-18.04/linkat-reassigna-equips-edc-18.04/debian/postinst
#!/bin/bash
# postinst script for linkat-reassigna-equips-edc
#
# see: dh_installdeb(1)

set -e

# summary of how this script can be called:
#        * <postinst> `configure' <most-recently-configured-version>
#        * <old-postinst> `abort-upgrade' <new version>
#        * <conflictor's-postinst> `abort-remove' `in-favour' <package>
#          <new-version>
#        * <postinst> `abort-remove'
#        * <deconfigured's-postinst> `abort-deconfigure' `in-favour'
#          <failed-install-package> <version> `removing'
#          <conflicting-package> <version>
# for details, see http://www.debian.org/doc/debian-policy/ or
# the debian-policy package

case "$1" in
    configure)
	if [ ! -e /etc/lightdm/lightdm.conf ]; then
		printf  "[SeatDefaults]\n\
greeter-show-manual-login=true\n" > /etc/lightdm/lightdm.conf
	fi
	#
	# Si existeix la wifi educaendigital, es canvia l'atribut del fitxer i s'elimina
	#
	if [ -e /etc/NetworkManager/system-connections/educaendigital ]; then
		chattr -i /etc/NetworkManager/system-connections/educaendigital
		rm /etc/NetworkManager/system-connections/educaendigital
	fi
	#
	# Es modifica la configuració de les xarxes WiFi de l'equip per compartir-la entre els usuaris de l'ordinador.
	# NetworkManager no mostra la contrasenya de la wifi quan es comparteix entre usuaris.
	# Es deixa el NetworkManager sota el control del PolicyKit 
	#
	if [ -d /etc/NetworkManager/system-connections ]; then
		find /etc/NetworkManager/system-connections -type f -exec sed -i '/permissions/c\permissions=' '{}' \;
		deb-systemd-invoke reload NetworkManager >/dev/null 2>&1
	fi
	#
	cp /usr/share/linkat/linkat-reassigna-equips-edc/linkat-reassigna-equip-edc-cron /etc/cron.d/
	chown root:root /etc/cron.d/linkat-reassigna-equip-edc-cron
	chmod 644 /etc/cron.d/linkat-reassigna-equip-edc-cron
	cp /usr/share/linkat/linkat-reassigna-equips-edc/linkat-reassigna-equips /etc/sudoers.d/
	chown root:root /etc/sudoers.d/linkat-reassigna-equips
	chmod 440 /etc/sudoers.d/linkat-reassigna-equips
	chmod 750 /usr/share/linkat/linkat-reassigna-equips-edc/*.sh
	chmod 755 /usr/share/linkat/linkat-reassigna-equips-edc/linkat-autoservei-usuari.sh
	chown root:root /usr/share/linkat/linkat-reassigna-equips-edc/*.sh
	invoke-rc.d cron reload >/dev/null 2>&1
	if [ -e /usr/share/linkat/linkat-reassigna-equips-edc/linkat-reassigna-equips-edc.desktop ]; then
		desktop-file-install /usr/share/linkat/linkat-reassigna-equips-edc/linkat-reassigna-equips-edc.desktop
	fi
	SUPORT_HOME=$(getent passwd suport |cut -d ":" -f 6)
	source "$SUPORT_HOME"/.config/user-dirs.dirs
	SUPORT_DESKTOP="$(sudo -H -u suport bash -c 'xdg-user-dir DESKTOP')"
	if [ -e "$SUPORT_HOME" ]; then
		if [ -e "$SUPORT_DESKTOP"/educa.desktop ]; then
			chattr -i "$SUPORT_DESKTOP"/educa.desktop
			chown -R suport:suport "$SUPORT_HOME/.dbus"
		fi
	fi
	sudo -H -u suport bash -c 'xdg-desktop-icon install /usr/share/linkat/linkat-reassigna-equips-edc/linkat-reassigna-equips-edc.desktop'
	chown suport:suport "$SUPORT_DESKTOP/linkat-reassigna-equips-edc.desktop"
	chmod 755 "$SUPORT_DESKTOP/linkat-reassigna-equips-edc.desktop"
	sudo -H -u suport bash -c "dbus-launch gio set $SUPORT_DESKTOP/linkat-reassigna-equips-edc.desktop metadata::trusted yes" >/dev/null 2>&1
    ;;

    abort-upgrade|abort-remove|abort-deconfigure)
    ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0
