# linkat-stat/src/etc/cron.d/linkat-stat-cron
SHELL=/bin/bash
@reboot root /usr/share/linkat/linkat-stat/linkat-stat.sh
