# linkat-stat/src/etc/cron.d/linkat-stat
SHELL=/bin/bash
@reboot root /usr/local/bin/linkat-stat.sh
