# linkat-stat/src/etc/cron.hourly/linkat-stat-cron-hourly
#!/bin/bash
/usr/share/linkat/linkat-stat/linkat-stat.sh HOURLY 1800
